# lib/system-skeleton/etc/inittab
# inittab # 

null::sysinit:/bin/mount -t proc proc /proc 
null::sysinit:/bin/mount -o remount,rw / 
null::sysinit:/bin/mkdir -p /dev/pts 
null::sysinit:/bin/mkdir -p /dev/shm 
null::sysinit:/bin/mount -a 
null::sysinit:/bin/hostname -F /etc/hostname 

tty1::respawn:/sbin/getty -L  tty1 115200 vt100 

null::shutdown:/bin/umount -a -r 
null::shutdown:/sbin/swapoff -a

